# crontab
* * * * * echo `date` >> /tmp/cron-test-log.txt
* * * * * node app help
